# extra/media.mount
[Unit]
Description=mount device to /media

[Mount]
What=/dev/sdc3
Where=/media

[Install]
WantedBy=basic.target
