# extra/mount_vboxshared.service
[Unit]
Description=Mount vbox shared dir

[Service]
ExecStart=/bin/mount /dev/sdb3 /media

[Install]
WantedBy=multi-user.target
