# extra/mount_vboxshared.service
[Unit]
Description=Mount vbox shared dir

[Service]
ExecStart=/bin/mount -t vboxsf tad /media

[Install]
WantedBy=multi-user.target
