# extra/mpd.service
[Unit]
Description=run mpd

[Service]
ExecStart=/usr/bin/mpd

[Install]
WantedBy=multi-user.target
